# packaging/udev/60-sovyx-audio-power.rules
# Sovyx — keep PCI audio-class devices at runtime_pm=on across
# hotplug (dock attach/detach), sleep/wake, and codec driver reload.
# Narrow target: ACTION add|change on SUBSYSTEM=sound, audio-class.
#
# Why this rule exists: the systemd oneshot fires only at boot. When
# the user docks/undocks a laptop, hot-swaps a USB-C hub that exposes
# an HDA audio device, or the kernel rebinds the codec driver (e.g.
# after a sleep/wake cycle), the child PCI device's power/control
# reverts to the kernel default ("auto"). Without re-asserting "on"
# on each such event, the pilot-class attenuation symptom reappears.
#
# Security model (updated after paranoid-QA audit — CRITICAL #4/#5):
# Earlier revisions of this rule used ``/bin/sh -c '… %S%p …'`` which
# string-interpolated the kernel-sourced sysfs path into a shell
# command. That path is normally clean, but (a) some USB descriptors
# can smuggle bytes that udev passes through verbatim and (b) future
# edits to the ATTRS filter could relax the sound-subsystem gate.
# The inline shell form made root-RCE possible once either of those
# assumptions broke.
#
# Current form: invoke ``/usr/libexec/sovyx/sovyx-audio-power-setone``
# with ``%S%p`` as a single argv element. The helper (POSIX sh)
# validates the argument — refusing any path not under /sys, any
# symlink that escapes /sys, and any target that isn't literally
# ``*/power/control``. Argv boundary means shell metacharacters in
# the path can never be reinterpreted.
#
# See ADR-voice-mixer-sanity-l2.5-bidirectional §H.2 + V2 Master
# Plan Part H + packaging/udev/sovyx-audio-power-setone.

ACTION=="change", SUBSYSTEM=="sound", \
    ATTRS{class}=="0x040300", \
    RUN+="/usr/libexec/sovyx/sovyx-audio-power-setone %S%p"

# Re-assert on add (new codec driver bind) so a runtime driver reload
# (e.g. modprobe -r snd_hda_intel && modprobe snd_hda_intel) restores
# our preferred state without a reboot.
ACTION=="add", SUBSYSTEM=="sound", \
    ATTRS{class}=="0x040300", \
    RUN+="/usr/libexec/sovyx/sovyx-audio-power-setone %S%p"
